# 推荐项目/04-泛娱乐推荐系统项目讲义&代码/recreation_project/supervisord.conf.d/start_redis.ini
[program:redis]
command=redis-server /root/recreation_project/conf/redis.conf
stdout_logfile=/root/recreation_project/log/redis_out.log
stderr_logfile=/root/recreation_project/log/redis_err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
